# examples/HaeorumAISearch/deploy/systemd/haeorum-ai-reindex.timer
[Unit]
Description=Run Haeorum Gift AI Search full reindex nightly

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=haeorum-ai-reindex.service

[Install]
WantedBy=timers.target
